# 99-monitor-hotplug.rules
#ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/start-dekstop.sh"
